# files/packages/lime-wireguard/root/etc/uci-defaults/90_wg-enable
!#/bin/sh

touch /etc/config/wireguard

[ -x /etc/init.d/wireguard ] &&
	/etc/init.d/wireguard enable
exit 0
